# contentloader.timer
[Unit]
Description=Run "contentloader --watch" twice daily

[Timer]
OnCalendar=*-*-* 00,12:00:00
Persistent=true

[Install]
WantedBy=timers.target
